# pkgs/dash/run.cron
#!/bin/bash
forever start -v -w --watchDirectory /home/dbond/menv/local/docs/menv/pkgs/dash/ index.js
