# pkgs/dash/run.cron
#!/bin/bash
forever start -v -w --watchDirectory /home/dbond/menv/pkgs/dash/ index.js
